# scripts/launchd/com.kilgore.deploy-freshness-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kilgore.deploy-freshness-watchdog</string>
    <!-- T-DEPLOY.5: detect persistent stale_alerted (manager behind origin/main)
         and auto-redeploy per manager-redeploy-gotchas-20260702.md. Runs OUTSIDE
         the fleet because fleet agents cannot restart their own manager (Gotcha 0). -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/kilgore/Dropbox/Code/cane/id-agents-deploy-main/scripts/deploy-freshness-watchdog.mjs</string>
    </array>
    <!-- every 15 minutes; acts only after 2 consecutive stale readings (>15 min). -->
    <key>StartInterval</key>
    <integer>900</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- git / npm / npx / launchctl / plutil must resolve; launchd's default
             PATH is minimal. Homebrew node first. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/kilgore</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/deploy-watchdog.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/deploy-watchdog.launchd.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
